# warrior-hq-resolve-geoip.timer
[Unit]
Description="Run Warrior HQ Resolve GeoIP hourly"

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
